# fortuna.service
[Unit]
Description=An AI powered daily quote generator!
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /usr/share/fortuna/fortuna.js

[Install]
WantedBy=default.target
