# fortuna.service
[Unit]
Description=An AI powered daily quote generator!

[Service]
Type=notify
ExecStart=fortuna

[Install]
WantedBy=multi-user.target